# enchtex-servers/web-server/files/templ/root/cron-scripts/mysql_optimized.cron
#!/bin/bash

mysqlcheck --repair --analyze --optimize --all-databases --auto-repair -u root -p
